# supervision/launchd/com.agentco.cc-health-check.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentco.cc-health-check</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{AGENT_CO_ROOT}}/agent-company/scripts/cc-version-check.sh</string>
        <string>check</string>
    </array>
    <key>WorkingDirectory</key>
    <string>{{AGENT_CO_ROOT}}</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>AGENT_CO_ROOT</key>
        <string>{{AGENT_CO_ROOT}}</string>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:{{HOME}}/.local/bin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>{{AGENT_CO_ROOT}}/logs/cc-health-check.log</string>
    <key>StandardErrorPath</key>
    <string>{{AGENT_CO_ROOT}}/logs/cc-health-check.log</string>
</dict>
</plist>
