# supervision/systemd/agentco-learning-consolidation.timer.template
[Unit]
Description=Agent Co — Learning Flywheel consolidation (daily 06:30)

[Timer]
OnCalendar=*-*-* 06:30:00
Persistent=true
AccuracySec=10min

[Install]
WantedBy=timers.target
