# supervision/systemd/agentco-cc-health-check.timer.template
[Unit]
Description=Agent Co — Daily Claude Code drift check at 07:00

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true
AccuracySec=10min

[Install]
WantedBy=timers.target
