# agent-company/scripts/templates/agentco-resource-monitor.timer.template
[Unit]
Description=Agent Company — Resource Monitor Timer (every 15 minutes)

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
AccuracySec=1min

[Install]
WantedBy=timers.target
