# agent-company/scripts/templates/agentco-resource-monitor.service.template
[Unit]
Description=Agent Company — Resource Monitor (cleanup stopped processes)

[Service]
Type=oneshot
ExecStart={{PROJECT_DIR}}/scripts/resource-monitor.sh
