# agent-company/scripts/templates/agentco-healthcheck.service.template
[Unit]
Description=Agent Company — Health Check (restart dead services)

[Service]
Type=oneshot
ExecStart={{PROJECT_DIR}}/scripts/healthcheck.sh
